# etc/crontabs/nginx
* * * * * source /etc/envvars_php; /usr/bin/php /var/www/localhost/htdocs/protected/yii queue/run 2>&1
* * * * * source /etc/envvars_php; /usr/bin/php /var/www/localhost/htdocs/protected/yii cron/run 2>&1
