# etc/service/cron/run
#!/bin/sh -e

exec crond -f -L /proc/self/fd/2